# nfc-listener.service
[Unit]
Description=Plexamp NFC Listener
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=andy
WorkingDirectory=/home/andy/Plexamp-NFC-Listener
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/andy/Plexamp-NFC-Listener/venv/bin/python -u /home/andy/Plexamp-NFC-Listener/nfc_listener.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
